# share/systemd/neural-tts.socket
[Unit]
Description=KDE TTS daemon synthesis socket
PartOf=neural-tts.service

[Socket]
ListenStream=%t/neural-tts.sock
SocketMode=0600
FileDescriptorName=neural-tts

[Install]
WantedBy=sockets.target
